# application/bin/vps-conciliation-masived.service
[Unit]
Description=Vial+ Platform Home and module integration management Service.
After=syslog.target

[Service]
ExecStart=/opt/vps/conciliation-masive/bin/conciliation-masive.jar
SuccessExitStatus=TEMPFAIL 250 SIGKILL

[Install]
WantedBy=multi-user.target